# warehouse.service
[Unit]
Description=Warehouse Management System
After=network.target postgresql.service

[Service]
User=warehouse
Group=warehouse
WorkingDirectory=/opt/warehouse
Environment="PATH=/opt/warehouse/venv/bin:/usr/local/bin:/usr/bin:/bin"
EnvironmentFile=/opt/warehouse/.env.production
ExecStart=/opt/warehouse/venv/bin/gunicorn --config gunicorn_config.py "app:create_app()"
Restart=always
RestartSec=5
StartLimitInterval=0

# Security hardening
PrivateTmp=true
ProtectSystem=full
NoNewPrivileges=true
PrivateDevices=true
ReadWritePaths=/opt/warehouse/instance /var/log/warehouse /var/backups/warehouse

[Install]
WantedBy=multi-user.target